# s6-overlay/s6-rc.d/svc-onedrive/finish
#!/bin/sh
if [ "$1" -ne 0 ]; then
    sleep 5
fi
