# s6-overlay/s6-rc.d/init-config/up
/command/with-contenv /scripts/entrypoint.sh
