# s6-overlay/s6-rc.d/svc-obsidian/run
#!/command/with-contenv sh
exec gosu obsidian /opt/obsidian/obsidian --no-sandbox --disable-gpu
